# supervisord.conf
[program:node]
command = make run
directory = /home/dotcloud/current
